# config/includes.chroot/etc/systemd/system/multi-user.target.wants/network-init.service
[Unit]
Description=Initialise le réseaux

[Service]
Type=oneshot
ExecStart=/usr/bin/network-init

[Install]
WantedBy=multi-user.target
